# app-laptop/prey/files/prey-agent
#!/sbin/runscript

_runfile="/var/run/prey-agent.pid"
_lockfile="/var/run/prey-agent.lock"

start(){
    ebegin "Starting Prey Agent"
    /usr/lib/node_modules/prey/bin/prey && echo $! > "${_runfile}" && touch "${_lockfile}"
    eend $?
}

stop(){
    ebegin "Stopping Prey Agent"
    kill -HUP $(cat "${_runfile}") && rm "${_runfile}" "${_lockfile}"
    eend $?
}
