# app-laptop/prey/files/prey-agent
#!/sbin/openrc-run
#

ARCH=$(uname -m)
[[ $ARCH == "x86_64" ]] && SUFFIX="64" || SUFFIX="32"

command="/usr/lib${SUFFIX}/prey/bin/prey"
pidfile="/var/run/prey-agent.pid"
description="The Prey Agent Daemon"
start_stop_daemon_args-"-b -n --pidfile ${pidfile}"
command_user="prey:prey"
