# app-shells/shellinabox/files/shellinaboxd
#!/sbin/runscript

_runfile="/var/run/shellinaboxd.pid"
_lockfile="/var/lock/shellinaboxd.lock"

start(){
    ebegin "Starting Shell In A Box Daemon"
    /usr/bin/shellinaboxd -b -s '/:SSH:localhost' --localhost-only -c '/tmp' --background="${_runfile}" && touch "${_lockfile}"
    eend $?
}

stop(){
    ebegin "Stopping Shell In A Box Daemon"
    kill -HUP $(cat "${_runfile}") && rm "${_runfile}" "${_lockfile}"
    eend $?
}
