# sys-apps/anything-sync-daemon/files/asd
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

extra_commands="sync"

_runfile="/var/run/asd.pid"

depend() {
	need localmount
}

start() {
	ebegin "Starting Anything Sync Daemon"
	/usr/bin/anything-sync-daemon sync && touch "${_runfile}"
	eend $?
}

stop() {
	ebegin "Stop Anything Sync Daemon"
	/usr/bin/anything-sync-daemon unsync && rm "${_runfile}"
	eend $?
}

sync() {
	ebegin "Doing a user-requested Sync"
	[ -f "$_runfile" ] && /usr/bin/anything-sync-daemon sync
	eend $?
}
