# x11-libs/bamf/files/bamfdaemon
#!/sbin/openrc-run

command=/usr/libexec/bamf/bamfdaemon
pidfile=${pidfile-/var/run/bamfdaemon.pid}
description="The BAMF daemon"

start(){
    ebegin "Starting BAMF daemon"
    start-stop-daemon --quiet --background --start --exec $command \
        --make-pidfile --pidfile $pidfile -- ""
    eend $? "Failed to start BAMF daemon"
}

stop(){
    ebegin "Stopping BAMF daemon"
    start-stop-daemon --quiet --stop --pidfile $pidfile --signal QUIT
    eend $? "Failed to stop BAMF daemon"
}
